# etc/network/if-up.d/htpdate
#!/bin/bash
SERVERS="www.linux.org www.freebsd.org"
/usr/bin/htpdate -t -s -i $SERVERS

